# rc.local
# Sample entry for '/etc/rc.local'
/bin/bash /opt/myfire/myfire-rebuild.sh &
